# profile.d/postgresql.sh
export PATH="/home/work/pg_5432/bin:$PATH"
export PGPORT="5432"
